# supervisord.conf
; Runs backend and frontend as two managed processes. Backend listens on
; 8001 and serves everything under /api, this environment's ingress is
; expected to route /api/* to port 8001 and everything else to the
; frontend's port (3000 by default). Both ports are overridable via each
; service's .env (copy from .env.example).
[supervisord]
nodaemon=true

[program:backend]
directory=%(here)s/backend
command=sh -c "uvicorn app:app --host 0.0.0.0 --port ${PORT:-8001}"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
directory=%(here)s/frontend
command=sh ./start.sh
environment=PORT="3000"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
